# config/supervisord/supervisord.proxy.conf
[unix_http_server]
file=/dev/shm/supervisor.sock
chmod=0700

[supervisord]
nodaemon=false
logfile=/run/supervisord.log
loglevel=warn
pidfile=/run/supervisord.pid

[program:httpd]
process_name=%(program_name)s_%(process_num)02d
command=httpd -D FOREGROUND
autostart=true
autorestart=false
startretries=3
numprocs=1
startsecs=0