# config/supervisord/vocgame.conf
[program:php-fpm]
process_name=%(program_name)s_%(process_num)02d
command=php-fpm --nodaemonize
autostart=true
autorestart=true